# ios-deb/layout/Library/LaunchDaemons/com.iphone-farm.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.iphone-farm.server</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/iphone-farm-server</string>
        <string>--port</string>
        <string>7777</string>
        <string>--fps</string>
        <string>15</string>
        <string>--quality</string>
        <string>75</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/var/log/iphone-farm-server.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/iphone-farm-server.log</string>
</dict>
</plist>
